# myapp-deb/DEBIAN/prerm
#!/bin/bash
systemctl stop myapp
systemctl disable myapp

